# deploy/com.sid.kalshi-paper-trader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sid.kalshi-paper-trader</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sid/Claude/Projects/Kalshi MLB market analysis/kalshi-pipeline/target/release/paper-trader</string>
        <string>--config</string>
        <string>config.toml</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/sid/Claude/Projects/Kalshi MLB market analysis/kalshi-pipeline</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/sid/Claude/Projects/Kalshi MLB market analysis/kalshi-pipeline/paper_trader_live.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sid/Claude/Projects/Kalshi MLB market analysis/kalshi-pipeline/paper_trader_live.log</string>
</dict>
</plist>
